# funding-demo/deploy/funding-demo.service
[Unit]
Description=Funding Demo FastAPI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=funddemo
Group=funddemo
WorkingDirectory=/opt/funding-demo
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/funding-demo/.venv/bin/python -m uvicorn web_app_core:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
